# systemd/pru-loader.service
[Unit]
Description=PRU0 CAN timestamp firmware — configure pins, enable OCP, load firmware
After=local-fs.target
Before=can-sniffer.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/can_sniffer/scripts/setup_pru.sh

[Install]
WantedBy=multi-user.target
